# maa-scheduler.service
[Unit]
Description=MAA任务调度器服务
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/Task/MAA_Auto
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 main.py run --config task_config.ini
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=maa-scheduler

# 环境变量
Environment=PYTHONPATH=/Task/MAA_Auto
Environment=PYTHONUNBUFFERED=1

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/Task/MAA_Auto

[Install]
WantedBy=multi-user.target