# maa-web-manager.service
[Unit]
Description=MAA Web Manager - Web界面管理MAA自动化任务
After=network.target maa-auto.service
Wants=maa-auto.service

[Service]
Type=simple
User=root
WorkingDirectory=/Task/MAA_Auto
Environment=PYTHONPATH=/Task/MAA_Auto
ExecStart=/usr/bin/python3 /Task/MAA_Auto/web_manager.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/Task/MAA_Auto
ProtectHome=true

[Install]
WantedBy=multi-user.target